# .config/systemd/user/commit-docs.timer
[Unit]
Description=Daily auto-commit of org docs

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
