# .config/systemd/user/mbsync.service
[Unit]
Description=Sync mail (mbsync backstop for goimapnotify)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/mailsync
